# systemd/hendricks-user-scheduler.service
[Unit]
Description=HP Hendricks 8A4D user-space fan scheduler
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hendricks-user-scheduler --profile performance --sensors CPU
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
